# linux/udev/41-aic8800d80-modeswitch.rules
# AIC8800D80 "Pandora" clone adapter - automatic mode switch
# Triggers usb_modeswitch when the device is plugged in as 1111:1111
ACTION=="add", ATTRS{idVendor}=="1111", ATTRS{idProduct}=="1111", RUN+="/usr/sbin/usb_modeswitch -c /etc/usb_modeswitch.d/1111:1111"

# After mode switch, the device re-enumerates as a69c:8d83 (WiFi mode).
# If aic8800_fdrv doesn't auto-bind, add the VID:PID dynamically:
ACTION=="add", ATTRS{idVendor}=="a69c", ATTRS{idProduct}=="8d83", RUN+="/bin/sh -c 'echo a69c 8d83 > /sys/bus/usb/drivers/aic8800_fdrv/new_id 2>/dev/null || true'"
